# pgadmin/files/etc/supervisord.d/php-fpm.ini
[program:php-fpm]
command=/opt/remi/php71/root/sbin/php-fpm -F
priority=100
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

