# nginx_route/files/etc/supervisord.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority=100
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

