# asterisk_13lts/install/supervisor/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
priority=999
process_name = httpd
redirect_stderr=true