# asterisk_13lts/install/supervisor/mysqld.ini
[program:mysqld]
command=/usr/libexec/mysqld --basedir=/usr --datadir=/var/lib/mysql --plugin-dir=/usr/lib64/mysql/plugin --user=mysql --log-error=/var/log/mariadb/mariadb.log --pid-file=/var/run/mariadb/mariadb.pid --socket=/var/lib/mysql/mysql.sock
user=mysql
priority=10
autostart=true
autorestart=true