# asterisk_13lts/install/supervisor/php-fpm.ini
[program:php-fpm]
command=/opt/remi/php56/root/usr/sbin/php-fpm -F
priority=500
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
