# freeswitch/install/supervisor/dbus.ini
[program:dbus]
command=dbus-daemon --system --nofork --nopidfile
priority=1
autostart=true
autorestart=true
redirect_stderr=true