# freeswitch/install/supervisor/freeswitch.ini
[program:freeswitch]
command=/usr/bin/freeswitch -nf -conf /etc/freeswitch -log /var/log/freeswitch -db /var/lib/freeswitch/db
priority=50
user=freeswitch
startsecs=10
stopwaitsecs = 600
redirect_stderr = false
stdout_logfile=/var/log/freeswitch/console.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=1000
stderr_logfile=/var/log/freeswitch/console.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=1000