# freeswitch/install/supervisor/memcached.conf
[program:memcached]
command=/usr/bin/memcached -m 256 -c 1024 -t 4
autostart=true
autorestart=true
user=memcache
priority=50
redirect_stderr=true
