# freeswitch/install/supervisor/postgresql.ini
[program:postgresql]
command=/usr/pgsql-9.4/bin/postgres -D /var/lib/pgsql/9.4/data
user=postgres
autostart=true
autorestart=true