# freeswitch/install/supervisor/postgresql.ini
[program:postgresql]
command=/usr/lib/postgresql/9.4/bin/postgres -D /var/lib/postgresql/9.4/main
user=postgres
autostart=true
autorestart=true