# cloud_pgadmin/files/etc/supervisord.d/pgadmin.ini
[program:pgadmin]
command=/usr/bin/uwsgi --socket 127.0.0.1:9005 --uid apache --gid apache --enable-threads --thunder-lock --pidfile /var/run/pgadmin.pid --plugin python3 --master --buffer-size=32768 --chdir /usr/lib/python3.4/site-packages/pgadmin4 --wsgi-file /usr/lib/python3.4/site-packages/pgadmin4/pgAdmin4.wsgi
priority=100
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
